# 99-mysensors.rules
ACTION=="add", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="db42", RUN+="/usr/bin/dfu-util -e"
